# deployment/lianai-game1.service
[Unit]
Description=LianAI Game Backend - AI-1 (develop_copaw)
Documentation=https://github.com/qiuxirufengQAQ/lianai
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/lianai1/backend

# 启动命令 - 端口 8081
ExecStart=/usr/bin/java -jar /var/www/lianai1/backend/lianai-game-1.0.0.jar

# 重启策略
Restart=on-failure
RestartSec=10

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 超时设置
TimeoutStartSec=60
TimeoutStopSec=30

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lianai-game1

# 环境变量 - 关键！确保隔离
Environment="SERVER_PORT=8081"
Environment="DB_NAME=lianai_game1"
Environment="SPRING_PROFILES_ACTIVE=ai1"
Environment="JAVA_OPTS=-Xms512m -Xmx2g"

# PID 文件
PIDFile=/var/www/lianai1/app.pid

[Install]
WantedBy=multi-user.target
